# supervisord.conf
[program:websocktify]
directory = /home/dotcloud/current/noVNC
command = python ./utils/websockify --web ./ 8080 idax-macewindu.azva.dotcloud.net:48073
